# system_files/shared/usr/lib/systemd/system/rouge-mok-enroll.service
[Unit]
Description=Stage Rouge Secure Boot keys for MOK enrollment
ConditionSecurity=uefi-secureboot
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/rouge/stage-mok-keys
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
